# .config/autostart/DB0F27D08451990662C36797B21595984E638B96.google-chrome-service.desktop
[Desktop Entry]
Type=Application
Terminal=false
Exec=/opt/google/chrome/chrome --type=service
Name=Google Chrome Service
